# systemd/user/weather-update.service
[Unit]
Description=Update weather cache

[Service]
Type=oneshot
ExecStart=/home/kazuki/.config/weather/weather-update.sh

